# crontab
* * * * * cd /database && python /project/populate.py latest

